# deployments/docker/build/s6files/services.d/dsiem/run
#!/command/with-contenv sh

[ -z "$DSIEM_NODE" ] && echo "DSIEM_NODE env variable must be set" && exit 1

if [ "$DSIEM_MODE" = "cluster-frontend" ]; then
  exec /dsiem/dsiem-frontend serve -n ${DSIEM_NODE} 2>&1
elif [ "$DSIEM_MODE" = "cluster-backend" ]; then
  exec /dsiem/dsiem-backend serve -n ${DSIEM_NODE} 2>&1
else
  echo "DSIEM_MODE env variable must be set to either 'cluster-frontend' or 'cluster-backend'"
  exit 1
fi
